# fuse/fuseesbdemo/files/supervisor.ini
[program:tomcat]
command=/opt/jboss/jboss-ews-2.0/tomcat7/bin/startup.sh
user=jboss
redirect_stderr=true
