# jon/jon-server/files/supervisor.ini
[program:postgresql]
command=scl enable postgresql92 './startpostgres.sh' 
redirect_stderr=true
priority=1

[program:jon]
command=/opt/jon/jon-server-3.2.0.GA/bin/rhqctl install --start
user=jboss
redirect_stderr=true
priority=2
autostart=false

[eventlistener:startjon]
command=python startjon.py
events=PROCESS_STATE_RUNNING
numprocs=1
