# jon/jon-agent/files/supervisor.ini
[program:agent]
command=/opt/jboss/rhq-agent/bin/rhq-agent-wrapper.sh start
user=jboss
redirect_stderr=true
