# eap/files/supervisor.ini
[program:eap]
command=/launch.sh
redirect_stderr=true
